# supervisor.conf
[program:strangerlights]
command=pipenv run python /home/pi/repo/strangerlights/strangerlights.py
directory=/home/pi/repo/strangerlights
user=root
autostart=true
autorestart=true
